# usr/lib/systemd/system/openrgb.service
[Unit]
Description=OpenRGB control

[Service]
Type=oneshot
ExecStart=/usr/bin/openrgb --server --verbose

[Install]
WantedBy=multi-user.target
